# Crawling/scrapy/etc/systemd/system/snap.certbot.renew.timer
[Unit]
# Auto-generated, DO NOT EDIT
Description=Timer renew for snap application certbot.renew
Requires=var-lib-snapd-snap-certbot-2618.mount
After=var-lib-snapd-snap-certbot-2618.mount
X-Snappy=yes

[Timer]
Unit=snap.certbot.renew.service
OnCalendar=*-*-* 02:56
OnCalendar=*-*-* 15:46

[Install]
WantedBy=timers.target
